# archive/dockerfiles/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/bin/bash -c "redis-server"

[program:consul]
command=/bin/bash -c "consul agent -config-file /consul/config/agent.conf $CONSUL_OPT"
